# debian/curvedns.service
[Unit]
Description=curvedns daemon
After=network.target

[Service]
Type=forking
EnvironmentFile=/etc/default/curvedns
ExecStart=/usr/sbin/curvedns $IP $PORT $REMOTEIP $REMOTEPORT
Restart=on-failure

[Install]
WantedBy=multi-user.target
